# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true

[program:api]
command=/usr/local/bin/python /usr/src/app/apl_api/api.py
autorestart=true